# conf/supervisor/celery.conf
[program:celery]
command=celery worker --workdir=/opt -A alertsd --time-limit=300 --concurrency=8 --config=celeryconfig
environment=C_FORCE_ROOT=1
autostart=true
autorestart=true
stdout_logfile          = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile          = /dev/stderr
stderr_logfile_maxbytes = 0
